# mysql/mysqld
#! /bin/bash
#
# mysqld   Start/Stop the MySQL mysqld Server
#
# chkconfig: 2345 90 60
# description: mysqld is the MySQL mysqld Server
# processname: mysqld
# config: /etc/mysqld_config
# pidfile: /var/run/mysqld.pid

user=mysqld
group=SYSTEM
PATH=/sbin:/usr/sbin:/usr/bin
PREFIX=/usr
CONFIG=/etc/mysqld_config
PIDFILE=/var/run/mysql/mysqld.pid
LOGFILE=/var/log/mysqld.log

# Source function library.
if [ -f /etc/init.d/functions ] ; then
  . /etc/init.d/functions
fi

RETVAL=0

# See how we were called.

prog="mysqld"
progdir="/usr/sbin"
program="$progdir/$prog.exe"

test -f $program || exit 0

# Source configuration
if [ -f /etc/sysconfig/$prog ] ; then
    . /etc/sysconfig/$prog
fi

start() {
	echo -n $"Starting $prog: "
	# check if cygrunsrv process
	cygrunsrv --start $prog
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch $PIDFILE  && echo "done."
	return $RETVAL
}

stop() {
	echo -n $"Stopping $prog: "
	cygrunsrv --stop $prog
	if test -r $PIDFILE; then
	    kill `cat $PIDFILE` 2> /dev/null || echo -n " failed"
	else
	    echo -n " no PID file"
	fi
	# this really needs a long time to stop.
	RETVAL=$?
	echo "."
	[ $RETVAL -eq 0 ] && rm -f $PIDFILE && echo "done."
	return $RETVAL
}

rhstatus() {
	# service status
	cygrunsrv -Q $prog
}

restart() {
	echo -n $"Restarting $prog: "
	$0 stop
	sleep 1
	$0 start
	echo "done."
}

install() {
	echo -n $"Installing $prog program: "
	# some safety measures
	touch $LOGFILE $PIDFILE
	chgrp $group $LOGFILE $PIDFILE
	chmod g+w $LOGFILE $PIDFILE
	
	cygrunsrv -I ${prog} -d "CYGWIN ${prog}" -p "${program}" -u "${user}" -e "CYGWIN=binmode tty" -a "--user mysqld --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib/mysql/plugins --log-error=/var/log/mysql/mysqld.err --pid-file=/var/run/mysql/mysqld.pid"
	
}
uninstall() {
	echo -n $"Uninstalling $prog program: "
	stop
	cygrunsrv --remove $prog
	echo "done."
}
reload() {
	echo -n $"Reloading $prog daemon configuration: "
	echo "unsupported. $0 restart"
	return 1

	/usr/bin/kill -HUP `cat $PIDFILE`
	RETVAL=$?
	[ $RETVAL -eq 0 ] && echo "done."
	RETVAL=$?
	return $RETVAL
}

case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  restart)
	restart
	;;
  reload)
	reload
	;;
  install)
	install
	;;
  uninstall)
	uninstall
	;;
  status)
	rhstatus
	;;
  condrestart)
	[ -f $PIDFILE ] && restart || :
	;;
  *)
	echo $"Usage: $0 {start|stop|status|install|uninstall|restart|condrestart}"
	exit 1
esac

exit $?
